# cookbooks/crawler/templates/default/etc/init/gunicorn.conf.erb
# This file is auto generated by chef
# Do no edit in place
description "Gunicorn for rdc_crawler"
author "Rafael Durán Castañeda <rafadurancastaneda@gmail.com>"
start on runlevel [2345]
stop on runlevel [!2345]

script
exec <%= node[:crawler][:venv_path] %>/bin/python <%= node[:crawler][:path] %>/rdc_crawler/manage.py run_gunicorn \
    --workers=<%= node[:crawler][:gunicorn_workers] %> --worker-class=<%= node[:crawler][:gunicorn_worker_class] %> --log-level=<%= node[:crawler][:log_level] %> --user=<%= node[:crawler][:user] %> \
    --access-logfile=<%= node[:gunicorn][:log_dir] %>/access.log \
    --error=<%= node[:gunicorn][:log_dir] %>/error.log
end script

respawn
